# D_systemd_autoupdate_service.timer
You are a senior Arch Linux developer. Output ONLY code blocks with file paths. No explanations unless critical.

[Unit]
Description=NeOS Auto Update Service
Wants=neos-autoupdate.timer

[Service]
Type=oneshot
ExecStart=/usr/local/bin/neos-autoupdate.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target