# profile/airootfs/etc/systemd/system/neos-autoupdate.timer
[Unit]
Description=NeOS Autoupdate Timer

[Timer]
OnBootSec=15min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
