# profile/airootfs/etc/systemd/system/neos-liveuser-setup.service
[Unit]
Description=NeOS Liveuser Setup
Before=display-manager.service sddm.service

[Service]
# Bolt: Monitor startup time overhead
# Palette: Ensure journal logs are clear
Type=oneshot
Environment=LC_ALL=C
ExecStart=/usr/local/bin/neos-liveuser-setup
# DO NOT sandbox this unit. Its whole job is system-level account creation:
# `useradd -m` needs lock files in /etc and must create /home/liveuser, and
# `sudo -u liveuser` needs new privileges. ProtectSystem=strict / ProtectHome=yes
# / NoNewPrivileges make those fail, the script dies on `set -e`, liveuser is
# never created, and SDDM drops to a manual greeter with no user (no autologin).
# An upstream merge (bc64daa) re-added that sandboxing and broke the live boot.

# Sentinel: Safe systemd sandboxing (omitting ProtectSystem/ProtectHome/NoNewPrivileges which break useradd)
ProtectKernelTunables=yes
ProtectControlGroups=yes
ProtectHostname=yes
RestrictRealtime=yes
LockPersonality=yes
ProtectKernelModules=yes
MemoryDenyWriteExecute=yes
RestrictNamespaces=yes

[Install]
WantedBy=multi-user.target
