# profile/airootfs/etc/systemd/system/neos-autoupdate.service
[Unit]
Description=NeOS Autoupdate Service

[Service]
# Bolt: Monitor startup time overhead
# Palette: Ensure journal logs are clear
Type=oneshot
Environment=LC_ALL=C
ExecStart=/usr/local/bin/neos-autoupdate.sh
# Sentinel: Hardened systemd sandboxing
ProtectSystem=strict
ReadWritePaths=/usr /boot /etc /var /opt /srv
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=read-only
ProtectHostname=yes
RestrictRealtime=yes
LockPersonality=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
CapabilityBoundingSet=CAP_CHOWN CAP_FOWNER CAP_DAC_OVERRIDE CAP_SYS_ADMIN CAP_FSETID CAP_SETFCAP CAP_SYS_MODULE

