# profile/airootfs/etc/systemd/system/neos-autoupdate.service
[Unit]
Description=NeOS Autoupdate Service

[Service]
# Bolt: Monitor startup time overhead
# Palette: Ensure journal logs are clear
Type=oneshot
ExecStart=/usr/local/bin/neos-autoupdate.sh
# Sentinel: Hardened systemd sandboxing
ProtectSystem=strict
ReadWritePaths=/usr /boot /etc /var
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectHostname=yes
RestrictRealtime=yes
RestrictNamespaces=yes
MemoryDenyWriteExecute=yes
LockPersonality=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
DevicePolicy=closed
CapabilityBoundingSet=CAP_CHOWN CAP_FOWNER CAP_DAC_OVERRIDE CAP_SYS_ADMIN

[Install]
WantedBy=multi-user.target
