# profile/airootfs/etc/skel/.config/autostart/neos-desktop-setup.desktop
[Desktop Entry]
Type=Application
Name=NeOS Desktop Setup
Exec=/usr/local/bin/neos-desktop-setup
X-KDE-autostart-phase=1
Hidden=false
NoDisplay=true
