# airootfs/etc/systemd/system/neos-autoupdate.service
[Unit]
Description=NeOS Automatic System Updates
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/neos-autoupdate.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
