# airootfs/etc/systemd/system/neos-driver-manager.service
[Unit]
Description=NeOS Driver Manager
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/neos-driver-manager
# Bolt: Verify that sandboxing additions do not measurably slow down service boot times.
# Palette: Ensure any failures related to tightened permissions yield clear, actionable log outputs for developers.
# Sentinel: Audit systemd sandboxing rules to prevent privilege escalation or unauthorized file system access.
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=yes
RestrictRealtime=yes

[Install]
WantedBy=multi-user.target
