# airootfs/etc/systemd/system/neos-performance-tweaks.service
[Unit]
Description=NeOS Performance Optimizations
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/neos-performance-tweaks
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target