# airootfs/etc/systemd/system/neos-driver-manager.service
[Unit]
Description=NeOS Driver Manager
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/neos-driver-manager
# Sentinel: Implement strict systemd service sandboxing (e.g., ProtectSystem=strict, NoNewPrivileges=yes, PrivateTmp=yes) to limit privileges.
ProtectSystem=full
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
