# airootfs/etc/systemd/system/neos-autoupdate.timer
[Unit]
Description=Timer for NeOS Automatic System Updates
Requires=neos-autoupdate.service

[Timer]
OnCalendar=weekly
OnBootSec=30min
Persistent=true

[Install]
WantedBy=timers.target
