# airootfs/etc/systemd/system/neos-autoupdate.service
[Unit]
Description=NeOS Autoupdate Service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/neos-autoupdate.sh
# Sentinel: Audit systemd sandboxing (ProtectSystem, NoNewPrivileges, PrivateTmp) to limit root privileges.
ProtectSystem=strict
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
