# airootfs/etc/systemd/system/neos-autoupdate.service
[Unit]
Description=NeOS Autoupdate Service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/neos-autoupdate.sh
# Sentinel: Implement strict systemd service sandboxing (e.g., ProtectSystem=strict, NoNewPrivileges=yes, PrivateTmp=yes) to limit privileges.

[Install]
WantedBy=multi-user.target
