# profile/airootfs/etc/systemd/system/sddm.service.d/10-software-greeter.conf
# Render the SDDM greeter with Qt Quick's software backend instead of OpenGL.
#
# In virtual machines without working 3D acceleration (notably VMware with
# "Accelerate 3D graphics" off, and some other guests), the GL-based greeter
# draws nothing — you get a black screen with only a mouse cursor. The software
# backend always renders. The cost is negligible: this only affects the static
# login screen, not the desktop session, so it is safe to leave on everywhere
# (bare metal included). The greeter inherits this from sddm.service's env.
[Service]
Environment=QT_QUICK_BACKEND=software
